# adnix.cron
MAILTO=root
PATH=/bin:/usr/bin:/usr/local/adnix

0 4 * * *  root  /usr/local/adnix/update_unbound
