# adnix.service
[Unit]
Description=Adnix updater

[Service]
ExecStart=/usr/local/adnix/update_unbound
